# deploy/launchd/com.neurostore.node.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neurostore.node</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/neuro-node</string>
        <string>--storage-path</string>
        <string>/Library/Application Support/NeuroStore/data</string>
        <string>--max-gb</string>
        <string>50</string>
        <string>--relay-url</string>
        <string>wss://demo.neurostore.network/v1/nodes/ws</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/neurostore-node.err</string>
    <key>StandardOutPath</key>
    <string>/var/log/neurostore-node.log</string>
</dict>
</plist>
